# omega_ai/mm_trap_detector/mm_trap_simulator.service
# MIT License
#
# Copyright (c) 2024 OMEGA BTC AI Team
#
# Permission is hereby granted, free of charge, to any person obtaining a copy
# of this software and associated documentation files (the "Software"), to deal
# in the Software without restriction, including without limitation the rights
# to use, copy, modify, merge, publish, distribute, sublicense, and/or sell
# copies of the Software, and to permit persons to whom the Software is
# furnished to do so, subject to the following conditions:
#
# The above copyright notice and this permission notice shall be included in all
# copies or substantial portions of the Software.
#
# THE SOFTWARE IS PROVIDED "AS IS", WITHOUT WARRANTY OF ANY KIND, EXPRESS OR
# IMPLIED, INCLUDING BUT NOT LIMITED TO THE WARRANTIES OF MERCHANTABILITY,
# FITNESS FOR A PARTICULAR PURPOSE AND NONINFRINGEMENT. IN NO EVENT SHALL THE
# AUTHORS OR COPYRIGHT HOLDERS BE LIABLE FOR ANY CLAIM, DAMAGES OR OTHER
# LIABILITY, WHETHER IN AN ACTION OF CONTRACT, TORT OR OTHERWISE, ARISING FROM,
# OUT OF OR IN CONNECTION WITH THE SOFTWARE OR THE USE OR OTHER DEALINGS IN THE
# SOFTWARE.

[Unit]
Description=BTC Market Maker Trap Simulator Service
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=omega_user
WorkingDirectory=/path/to/omega-btc-ai
ExecStart=/usr/bin/python3 -m omega_ai.mm_trap_detector.trap_simulation_service --volatility 1.2 --frequency 0.25 --sleep 0.1

# Environment
Environment="PYTHONUNBUFFERED=1"

# Restart settings
Restart=on-failure
RestartSec=5s

# Log settings
StandardOutput=journal
StandardError=journal

# Security settings
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target 