# divine_dashboard_v2/com.omega.divinebrowser.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.omega.divinebrowser.v2</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>REPLACE_WITH_FULL_PATH/divine_book_browser_v2/divine_server.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>REPLACE_WITH_FULL_PATH/divine_book_browser_v2</string>
    <key>StandardErrorPath</key>
    <string>REPLACE_WITH_FULL_PATH/divine_book_browser_v2/divine_server_error.log</string>
    <key>StandardOutPath</key>
    <string>REPLACE_WITH_FULL_PATH/divine_book_browser_v2/divine_server_out.log</string>
</dict>
</plist> 